# safetynet/layout/DEBIAN/postinst
#!/usr/bin/bash

launchctl load /Library/LaunchDaemons/com.nito.safetynet.plist